# package/kernel/aic8800/files/etc/init.d/aic8800-init
#!/bin/sh /etc/rc.common
# AIC8800 启动初始化 - 作为 hotplug 的 fallback

START=19

boot() {
    # 延迟等待 USB 和 wireless 就绪
    (
        sleep 15
        
        # 检查是否存在 AIC8800 phy 设备但未配置
        for phydir in /sys/class/ieee80211/phy*; do
            [ -d "$phydir" ] || continue
            
            devpath=$(readlink "$phydir/device" 2>/dev/null)
            [ -z "$devpath" ] && continue
            
            usbdev=$(basename "$devpath")
            vendor_path="/sys/bus/usb/devices/${usbdev}/idVendor"
            [ ! -f "$vendor_path" ] && vendor_path="/sys/bus/usb/devices/${usbdev%:*}/idVendor"
            [ ! -f "$vendor_path" ] && continue
            
            vendor=$(cat "$vendor_path" 2>/dev/null)
            [ "$vendor" = "a69c" ] || continue
            
            # 是 AIC8800，获取 phy 名
            phyname=$(basename "$phydir")
            
            # 检查是否已配置
            configured=0
            for radio in $(uci show wireless 2>/dev/null | grep '=wifi-device' | cut -d= -f1 | cut -d. -f2); do
                radiopath=$(uci -q get wireless.${radio}.path)
                if [ "$radiopath" = "$usbdev" ] || [ "$radiopath" = "usb1/1-1/${usbdev}" ] || [ "$radiopath" = "usb2/2-1/${usbdev}" ]; then
                    configured=1
                    # 修复路径
                    if [ "$radiopath" != "$usbdev" ]; then
                        uci set wireless.${radio}.path="$usbdev"
                        uci commit wireless
                        logger -t "aic8800-init" "Fixed ${radio} path: $radiopath -> $usbdev"
                    fi
                    break
                fi
            done
            
            if [ "$configured" = "0" ]; then
                # 创建新配置
                idx=0
                while uci -q get wireless.radio${idx} >/dev/null 2>&1; do
                    idx=$((idx + 1))
                done
                
                uci batch << EOF
set wireless.radio${idx}=wifi-device
set wireless.radio${idx}.type=mac80211
set wireless.radio${idx}.path='${usbdev}'
set wireless.radio${idx}.band=5g
set wireless.radio${idx}.channel=36
set wireless.radio${idx}.htmode=HE80
set wireless.radio${idx}.disabled=0
set wireless.default_radio${idx}=wifi-iface
set wireless.default_radio${idx}.device=radio${idx}
set wireless.default_radio${idx}.network=lan
set wireless.default_radio${idx}.mode=ap
set wireless.default_radio${idx}.ssid=ZagWIFI
set wireless.default_radio${idx}.encryption=none
commit wireless
EOF
                logger -t "aic8800-init" "Created radio${idx} for AIC8800 ${phyname}"
            fi
            
            # 启用无线
            /sbin/wifi up radio${idx} 2>/dev/null || /sbin/wifi reload
            logger -t "aic8800-init" "WiFi enabled for AIC8800 ${phyname}"
        done
    ) &
}
