# package/kernel/aic8800/files/etc/hotplug.d/ieee80211/20-aic8800
#!/bin/sh
# AIC8800 自动配置脚本 - phy 创建时实时识别并配置
# 不依赖固定路径，通过 USB Vendor ID 实时识别

[ "${ACTION}" = "add" ] || exit 0
[ -n "${DEVNAME}" ] || exit 0

# 获取 phy 设备路径
phydir="/sys/class/ieee80211/${DEVNAME}"
[ -d "$phydir" ] || exit 0

devpath=$(readlink "$phydir/device" 2>/dev/null)
[ -z "$devpath" ] && exit 0

# 获取设备名称并尝试读取 USB Vendor ID
usbdev=$(basename "$devpath")
vendor_path="/sys/bus/usb/devices/${usbdev}/idVendor"
if [ ! -f "$vendor_path" ]; then
    # 尝试上级设备（处理 interface 节点情况）
    vendor_path="/sys/bus/usb/devices/${usbdev%:*}/idVendor"
fi

# 检查是否是 USB 设备（idVendor 文件存在）
[ -f "$vendor_path" ] || exit 0

vendor=$(cat "$vendor_path" 2>/dev/null)
[ "$vendor" = "a69c" ] || exit 0

# 是 AIC8800，获取当前 path
uci_path="$usbdev"

logger -t "aic8800" "AIC8800 detected: ${DEVNAME}, USB path: ${uci_path}"

# 检查是否已有该 USB 设备的配置（匹配任何包含此路径的 radio）
existing_radio=""
for j in 0 1 2 3 4 5; do
    cfg_path=$(uci -q get wireless.radio${j}.path 2>/dev/null)
    if [ -n "$cfg_path" ]; then
        # 检查是否是同一个设备（path 相同或包含）
        case "$cfg_path" in
            "$uci_path"|*"$uci_path")
                existing_radio="$j"
                break
                ;;
        esac
    fi
done

# 如果没有现有配置，找到空闲 radio 索引
if [ -z "$existing_radio" ]; then
    idx=0
    while uci -q get wireless.radio${idx} >/dev/null 2>&1; do
        idx=$((idx + 1))
    done
    
    # 创建新配置
    uci batch << EOF
set wireless.radio${idx}=wifi-device
set wireless.radio${idx}.type=mac80211
set wireless.radio${idx}.path='${uci_path}'
set wireless.radio${idx}.band=5g
set wireless.radio${idx}.channel=36
set wireless.radio${idx}.htmode=HE80
set wireless.radio${idx}.disabled=0
set wireless.default_radio${idx}=wifi-iface
set wireless.default_radio${idx}.device=radio${idx}
set wireless.default_radio${idx}.network=lan
set wireless.default_radio${idx}.mode=ap
set wireless.default_radio${idx}.ssid=ZagWIFI
set wireless.default_radio${idx}.encryption=none
commit wireless
EOF
    logger -t "aic8800" "Created radio${idx} for AIC8800"
    existing_radio="$idx"
else
    # 修复可能的错误 path 格式
    cfg_path=$(uci -q get wireless.radio${existing_radio}.path 2>/dev/null)
    if [ "$cfg_path" != "$uci_path" ]; then
        uci set wireless.radio${existing_radio}.path="$uci_path"
        uci commit wireless
        logger -t "aic8800" "Fixed radio${existing_radio} path: $cfg_path -> $uci_path"
    fi
fi

# 延时启用无线
(
    sleep 2
    /sbin/wifi up radio${existing_radio} 2>/dev/null || /sbin/wifi reload
    logger -t "aic8800" "WiFi enabled for radio${existing_radio}"
) &

exit 0
